# packaging/kakaotalk-wayland-sni.service.in
[Unit]
Description=StatusNotifier companion for KakaoTalk on Wine Wayland

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/kakaotalk-wayland-sni
EnvironmentFile=-%h/.config/kakaotalk-wayland-sni/config
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
